# .config/systemd/user/auto-specder.service
[Unit]
Description=Auto-spec proposer based on metrics/handoffs

[Service]
Type=oneshot
ExecStart=/home/kloros/bin/auto_specder.py
WorkingDirectory=/home/kloros
