# .config/systemd/user/library-extract.service
[Unit]
Description=Extract pattern library from winners

[Service]
Type=oneshot
ExecStart=/home/kloros/bin/library_extract.py
WorkingDirectory=/home/kloros
