# .config/systemd/user/repairlab-queue-watcher.timer
[Unit]
Description=Start/keep RepairLab queue watcher alive

[Timer]
OnBootSec=5
OnUnitActiveSec=5

[Install]
WantedBy=timers.target
