# D-REAM_BUNDLE/systemd/phase-heuristics.service
[Unit]
Description=PHASE Heuristic Controller
Documentation=https://github.com/kloros/phase
After=network.target

[Service]
Type=oneshot
User=kloros
Group=kloros
WorkingDirectory=/home/kloros
Environment="PYTHONPATH=/home/kloros/src"
ExecStart=/home/kloros/.venv/bin/python /home/kloros/src/heuristics/controller.py
StandardOutput=append:/home/kloros/logs/heuristics_controller.log
StandardError=append:/home/kloros/logs/heuristics_controller.log
TimeoutSec=120

[Install]
WantedBy=multi-user.target
