# .config/systemd/user/repairlab-queue-watcher.service.d/override.conf
[Service]
Restart=on-failure
RestartSec=2s
ExecReload=/bin/kill -HUP $MAINPID
