# src/core/config/systemd/kloros-cognitive-actions.service
[Unit]
Description=KLoROS Cognitive Actions Subscriber (Tier 2 - MEMORY_PRESSURE)
After=network.target kloros-chem-proxy.service
Requires=kloros-chem-proxy.service

[Service]
Type=simple
User=kloros
WorkingDirectory=/home/kloros/src
Environment=PYTHONPATH=/home/kloros/src
ExecStart=/home/kloros/.venv/bin/python -u consciousness/cognitive_actions_subscriber.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
